# thinkfan/thinkfan.service.d/00-generate-config.conf
[Service]
ExecStartPre=-/usr/bin/thinkfan-config.sh
